# etc/systemd/system/socksfwd@.service
[Unit]
Description=tun -> socks5

[Service]
ExecStart=/home/ray/bin/socksfwd %I
