# home/.config/systemd/user/break.timer
[Unit]
Description=Break Encouragement

[Timer]
OnBootSec=60min
OnUnitActiveSec=60min

[Install]
WantedBy=timers.target
