# home/.config/systemd/user/break.timer
[Unit]
Description=Break Encouragement

[Timer]
OnBootSec=30min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
