# home/.config/systemd/user/graphical.target
[Unit]
Description=Graphical Interface
AllowIsolate=yes
