# home/.config/systemd/user/inotify-tmp-wallpaper.service
[Service]
Description=Set /tmp/wallpaper.jpg as wallpaper
Type=oneshot
ExecStart=/usr/bin/feh --bg-scale /tmp/wallpaper.jpg
