# home/.config/systemd/user/inotify-tmp-wallpaper.path
[Path]
PathModified=/tmp/wallpaper.jpg

[Install]
WantedBy=graphical.target
