# home/.config/systemd/user/inotify-tmp-screen.service
[Service]
Type=oneshot
ExecStart=/usr/bin/rsync -a /tmp/screen.jpg linode-ca:/tmp/
