# home/.config/systemd/user/inotify-tmp-screen.path
[Path]
PathModified=/tmp/screen.jpg

[Install]
WantedBy=graphical.target
