# home/.config/systemd/user/nvchecker.service
[Unit]
Description=Check for new versions

[Service]
ExecStart=/usr/bin/nvchecker -n %h/abs/nvchecker.ini
