# home/.config/systemd/user/nvchecker.timer
[Unit]
Description=Check for new versions

[Timer]
OnBootSec=10min
OnUnitActiveSec=180min

[Install]
WantedBy=timers.target
