# home/.config/systemd/user/nvchecker.service
[Unit]
Description=Check for new versions

[Service]
ExecStart=/usr/bin/nvchecker -n /home/ray/abs/nvchecker.ini
