# home/.config/systemd/user/redshift-gtk.service
[Unit]
Description=Redshift display colour temperature adjustment (GUI)
Documentation=http://jonls.dk/redshift/
After=display-manager.service

[Service]
ExecStart=/usr/bin/redshift-gtk
Restart=always

[Install]
WantedBy=default.target
