# home/.config/systemd/user/taskbar.service
[Unit]
Description=Taskbar

[Service]
ExecStart=/home/ray/bin/taskbar

[Install]
WantedBy=graphical.target
