# home/.config/systemd/user/xscreensaver.service
[Service]
ExecStart=/usr/bin/xscreensaver -no-splash

[Install]
WantedBy=graphical.target
